# ubt-core/src/sv.conf
; /etc/supervisor/supervisor.conf
; [supervisord]
; logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
; logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
; logfile_backups=10           ; (num of main logfile rotation backups;default 10)
; loglevel=info                ; (log level;default info; others: debug,warn,trace)
; pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
; nodaemon=false               ; (start in foreground if true;default false)
; minfds=1024                  ; (min. avail startup file descriptors;default 1024)
; minprocs=200                 ; (min. avail process descriptors;default 200)
; user=pi                 ; (default is current user, required if root)
; environment=DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus/system_bus_socket"     ; (key value pairs to add to environment)

# TODO: if DBUS_SESSION_BUS_ADDRESS set?
# [program:dbus]
# ; --address="unix:path=/run/dbus/system_bus_socket"
# command=/usr/bin/dbus-daemon --system --nofork
# startsecs=2
# autorestart=false
# startretries=1 
# stdout_logfile=/var/log/supervisor/_dbus.log
# stdout_logfile_maxbytes = 50MB
# stdout_logfile_backups  = 10
# redirect_stderr=true

[program:xrdp]
priority=31
user=root
; autostart=true
; autorestart=true
; stopsignal=QUIT
startretries=5
autorestart=true
command=/usr/sbin/xrdp -nodaemon
stdout_logfile=/var/log/supervisor/_xrdp.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
redirect_stderr=true

[program:ssh]
priority=33
user=root
startretries=5
autorestart=true
command=/usr/sbin/dropbear -EFRp 22
stdout_logfile=/var/log/supervisor/_ssh.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
redirect_stderr=true

[program:webvnc]
priority=34
user=root
startretries=5
autorestart=true
command=bash -c "cd /usr/local/webhookd; exec bash ./run.sh"
stdout_logfile=/var/log/supervisor/_novnc.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
redirect_stderr=true
